# alpine-php/rootfs/etc/services.d/php/run
#!/command/with-contenv sh

exec /usr/sbin/php-fpm81 ${PHP_OPTIONS}
