# alpine-php/rootfs/etc/services.d/php/run
#!/command/with-contenv sh

exec /usr/sbin/php-fpm7 ${PHP_OPTIONS}
