# alpine-unbound/rootfs/etc/services.d/unbound/run
#!/command/with-contenv sh

exec /usr/sbin/unbound ${UNBOUND_OPTIONS}
