# alpine-xvfb/rootfs/etc/services.d/xvfb/run
#!/usr/bin/with-contenv sh

exec /usr/bin/Xvfb -ac -listen tcp :99
