# alpine-ipfs/rootfs/etc/services.d/ipfs/run
#!/command/with-contenv sh

export USER=ipfs
export HOME=/var/lib/ipfs

exec s6-setuidgid ipfs /usr/bin/ipfs daemon ${IPFS_OPTIONS}
