# alpine-rspamd/rootfs/etc/services.d/rspamd/run
#!/command/with-contenv sh

cd /var/lib/rspamd || exit 1
exec /usr/sbin/rspamd -f -u rspamd -g rspamd ${RSPAMD_OPTIONS}
