# alpine-docker-registry/rootfs/etc/services.d/docker-registry/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid docker-registry /usr/bin/docker-registry serve /etc/docker-registry/config.yml
