# alpine-portainer/rootfs/etc/services.d/portainer/run
#!/usr/bin/with-contenv sh

cd /usr/local/portainer || exit 1
exec s6-setuidgid portainer /usr/local/portainer/portainer -d "${PORTAINER_HOME}" ${PORTAINER_OPTIONS}
