# alpine-portainer/rootfs/etc/services.d/portainer/run
#!/usr/bin/with-contenv sh

cd /usr/local/portainer || exit 1
exec s6-setuidgid portainer ./portainer -p :8080 -d ${PORTAINER_HOME}
