# alpine-portainer/rootfs/etc/services.d/portainer/run
#!/usr/bin/with-contenv sh

exec /usr/local/portainer/portainer -d "${PORTAINER_HOME}" ${PORTAINER_OPTIONS}
