# alpine-build-nodejs/root/etc/services.d/build/run
#!/usr/bin/with-contenv sh

cd /usr/local/build
exec /bin/sh ${BUILD_SCRIPT}
