# alpine-nginx/rootfs/etc/services.d/nginx/run
#!/command/with-contenv sh

exec /usr/sbin/nginx ${NGINX_OPTIONS}
