# alpine-node-red/rootfs/etc/services.d/node-red/run
#!/command/with-contenv sh

cd /usr/local/node-red || exit 2
exec /usr/bin/node red.js --userDir /var/lib/node-red/ ${NODERED_OPTIONS}
