# alpine-mariadb/rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv sh

exec /usr/bin/mysqld --user=mysql --console --skip-networking=0 ${MARIADB_OPTIONS}
