# alpine-base-consul/root/etc/services.d/consul-template/run
#!/usr/bin/with-contenv sh

exec consul-template -consul ${CONSUL_AGENT}:8500 -config /etc/consul-template/conf-services.d ${CONSULTEMPLATE_OPTIONS}
