# alpine-base-consul/root/etc/services.d/consul-template/run
#!/usr/bin/with-contenv sh

exec /bin/consul-template -consul ${CONSUL_AGENT}:${CONSUL_PORT} -config /etc/consul-template/conf-services.d ${CONSULTEMPLATE_OPTIONS}
