# alpine-rabbitmq/root/etc/services.d/rabbitmq/run
#!/usr/bin/with-contenv sh

exec /usr/local/rabbitmq/sbin/rabbitmq-server;
