# alpine-rabbitmq/root/etc/services.d/rabbitmq/run
#!/usr/bin/with-contenv sh

exec rabbitmq-server
