# alpine-tomcat/root/etc/services.d/tomcat/run
#!/usr/bin/with-contenv sh

cd /usr/local/tomcat
exec /bin/sh bin/catalina.sh run
