# alpine-tomcat/root/etc/services.d/tomcat/run
#!/usr/bin/with-contenv sh

exec /bin/sh /usr/local/tomcat/bin/catalina.sh run
