# alpine-backuppc/rootfs/etc/services.d/lighttpd/run
#!/command/with-contenv sh

exec /usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D
