# alpine-hostapd/rootfs/etc/services.d/hostapd/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/hostapd /etc/hostapd/hostapd.conf
