# alpine-n8n/rootfs/etc/services.d/n8n/run
#!/usr/bin/with-contenv sh

export USER=n8n
export HOME=/var/lib/n8n

cd /var/lib/n8n || exit 2
exec s6-setuidgid n8n /usr/bin/npx n8n ${N8N_OPTIONS}
