# alpine-telegraf/rootfs/etc/services.d/telegraf/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid telegraf /usr/local/telegraf/usr/bin/telegraf --config /etc/telegraf/telegraf.conf
