# alpine-telegraf/rootfs/etc/services.d/telegraf/run
#!/usr/bin/with-contenv sh

cd /usr/local/telegraf || exit 1
exec s6-setuidgid telegraf /usr/local/telegraf/usr/bin/telegraf --config etc/telegraf/telegraf.conf ${TELEGRAF_OPTIONS}
