# alpine-opendmarc/rootfs/etc/services.d/opendmarc/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/opendmarc -f -c /etc/opendmarc/opendmarc.conf -u opendmarc -P /var/run/opendmarc/opendmarc.pid \
    -p inet:8893
