# alpine-strongswan/rootfs/etc/services.d/strongswan-ping/run
#!/command/with-contenv sh

exec /usr/local/bin/strongswan-ping ${STRONGSWANPING_OPTIONS}
