# alpine-dnsmasq/rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/dnsmasq --no-daemon ${DNSMASQ_OPTIONS}
