# alpine-zabbix-agent/rootfs/etc/services.d/zabbix-docker/run
#!/usr/bin/with-contenv sh

exec /usr/bin/python3 -m zabbix-docker.main ${ZABBIXDOCKER_OPTIONS}
