# alpine-zabbix-agent/rootfs/etc/services.d/zabbix-docker/run
#!/usr/bin/with-contenv sh

exec /bin/python3 -m zabbix-docker.main
