# alpine-zabbix-agent/rootfs/etc/services.d/zabbix-docker/run
#!/usr/bin/with-contenv sh

exec /usr/local/zabbix-docker/bin/zabbix-docker
