# alpine-nginx/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

exec nginx ${NGINX_OPTIONS}
