# alpine-postgresql/root/etc/services.d/postgresql/run
#!/usr/bin/with-contenv sh

exec su -c "postgres ${POSTGRESQL_OPTIONS}" postgres
