# alpine-jetty/rootfs/etc/services.d/jetty/run
#!/usr/bin/with-contenv sh

cd "${JETTY_BASE}" || exit 1
exec /usr/bin/java ${JAVA_OPTS} -jar "${JETTY_HOME}/start.jar"
