# alpine-openldap/rootfs/etc/services.d/openldap/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/slapd -u ldap -g ldap -h "ldap:///" -d0
