# alpine-tomcat/rootfs/etc/services.d/tomcat/run
#!/usr/bin/with-contenv sh

cd /usr/local/tomcat || exit 1
exec /usr/local/tomcat/bin/catalina.sh run
