# alpine-geth/rootfs/etc/services.d/geth/finish
#!/usr/bin/execlineb -S1

if { eltest ${1} -ne 0 }
if { eltest ${1} -ne 256 }

s6-svscanctl -t /run/service
