# alpine-rng/rootfs/etc/services.d/rng/run
#!/usr/bin/with-contenv sh

exec rngd -f ${RNG_OPTIONS}
