# alpine-rng/rootfs/etc/services.d/rng/run
#!/command/with-contenv sh

exec rngd -f ${RNG_OPTIONS}
