# alpine-dovecot/rootfs/etc/services.d/dovecot/run
#!/command/with-contenv sh

exec /usr/sbin/dovecot -F ${DOVECOT_OPTIONS}
