# alpine-dovecot/rootfs/etc/services.d/dovecot/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/dovecot -F ${DOVECOT_OPTIONS}
