# alpine-kibana/rootfs/etc/services.d/kibana/run
#!/usr/bin/with-contenv sh

cd /usr/local/kibana || exit 1
exec s6-setuidgid kibana /usr/local/kibana/bin/kibana
