# alpine-rundeck/rootfs/etc/services.d/rundeck/run
#!/command/with-contenv sh

cd ${RDECK_BASE} || exit 1
exec s6-setuidgid rundeck /usr/bin/java ${JAVA_OPTS} -jar rundeck.war --skipinstall {RUNDECK_OPTIONS}
