# alpine-rundeck/rootfs/etc/services.d/rundeck/run
#!/usr/bin/with-contenv sh

cd ${RDECK_BASE} || exit 1
exec s6-setuidgid rundeck /usr/bin/java ${JAVA_OPTS} \
    -jar rundeck.jar --skipinstall {RUNDECK_OPTS}
