# alpine-haproxy/rootfs/etc/services.d/haproxy-1/run
#!/usr/bin/with-contenv sh

old_pid=$(pidof haproxy)
echo $$ > /var/run/haproxy.pid

exec /usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg -db -st ${old_pid} ${HAPROXY_OPTIONS}
