# alpine-base-consul/root/etc/services.d/consul/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid consul consul agent -config-dir /etc/consul.d -retry-join ${CONSUL_SERVER} -retry-interval 15s
