# alpine-traefik/rootfs/etc/services.d/traefik/run
#!/command/with-contenv sh

exec /usr/local/traefik/traefik ${TRAEFIK_OPTIONS}
