# alpine-rsyslog/rootfs/etc/services.d/rsyslog/run
#!/command/with-contenv sh

exec /usr/bin/rsyslogd -n ${RSYSLOG_OPTIONS}
