# alpine-strongswan/rootfs/etc/services.d/strongswan-monitor/run
#!/command/with-contenv sh

exec /usr/local/bin/strongswan-monitor ${STRONGSWANMONITOR_OPTIONS}
