# alpine-strongswan/rootfs/etc/services.d/strongswan-monitor/run
#!/usr/bin/with-contenv sh

exec /usr/local/bin/strongswan-monitor
