# alpine-jetty/root/etc/services.d/jetty/run
#!/usr/bin/with-contenv sh

cd /usr/local/jetty
exec /usr/bin/java -jar start.jar
