# alpine-vault/rootfs/etc/services.d/vault/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid vault /bin/vault server -config /etc/vault.d ${VAULT_OPTIONS}
