# alpine-backuppc/rootfs/etc/services.d/spawn-fcgi/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid backuppc /usr/bin/spawn-fcgi -d /var/www/htdocs -p 9000 -n /usr/bin/fcgiwrap
