# alpine-mongodb/rootfs/etc/services.d/mongodb/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid mongodb /usr/bin/mongod -f /etc/mongod.conf ${MONGODB_OPTIONS}
