# alpine-rslsync/rootfs/etc/services.d/rslsync/run
#!/usr/bin/with-contenv sh

exec /usr/local/bin/rslsync --nodaemon --config /etc/rslsync.conf
