# alpine-rslsync/rootfs/etc/services.d/rslsync/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid rslsync /usr/local/bin/rslsync --nodaemon --config /etc/rslsync/rslsync.conf
