# alpine-influxdb/rootfs/etc/services.d/influxdb/run
#!/usr/bin/with-contenv sh

cd /usr/local/influxdb || exit 1
exec s6-setuidgid influxdb /usr/local/influxdb/usr/bin/influxd run -config ./etc/influxdb/influxdb.conf \
    ${INFLUXDB_OPTIONS}
