# alpine-backuppc/rootfs/etc/services.d/backuppc/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid backuppc /usr/share/BackupPC/bin/BackupPC
