# alpine-strongswan/rootfs/etc/services.d/strongswan/run
#!/command/with-contenv sh

exec /usr/sbin/ipsec start --nofork ${STRONGSWAN_OPTIONS}
