# alpine-strongswan/rootfs/etc/services.d/strongswan/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/ipsec start --nofork
