# alpine-kapacitor/rootfs/etc/services.d/kapacitor/run
#!/command/with-contenv sh

cd /usr/local/kapacitor || exit 1
exec s6-setuidgid kapacitor /usr/local/kapacitor/usr/bin/kapacitord -config ./etc/kapacitor/kapacitor.conf \
    ${KAPACITOR_OPTIONS}
