# alpine-kapacitor/rootfs/etc/services.d/kapacitor/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid kapacitor /usr/local/kapacitor/usr/bin/kapacitord -config /etc/kapacitor/kapacitor.conf
