# alpine-ndppd/rootfs/etc/services.d/ndppd/run
#!/command/with-contenv sh

exec /usr/local/sbin/ndppd -c /etc/ndppd.conf -v ${NDPPD_OPTIONS}
