# alpine-mailhog/rootfs/etc/services.d/mailhog/run
#!/command/with-contenv sh

exec s6-setuidgid mailhog /usr/local/bin/MailHog ${MAILHOG_OPTIONS}
