# alpine-zabbix-agent/rootfs/etc/services.d/zabbix-agent/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid zabbix /usr/sbin/zabbix_agentd -f
