# alpine-postfix/rootfs/etc/services.d/postfix/run
#!/usr/bin/with-contenv sh

exec /usr/bin/nodaemon /var/spool/postfix/pid/master.pid /usr/sbin/postfix start
