# alpine-postfix/rootfs/etc/services.d/postfix/run
#!/usr/bin/with-contenv sh

/etc/init.d/postfix start
/etc/init.d/postfix stop

exec /usr/lib/postfix/master -c /etc/postfix -d 2>&1
