# alpine-postfix/rootfs/etc/services.d/postfix/run
#!/usr/bin/with-contenv sh

exec /usr/lib/postfix/master -c /etc/postfix -d 2>&1
