# alpine-postfix/rootfs/etc/services.d/postfix/run
#!/command/with-contenv sh

exec /usr/libexec/postfix/master -c /etc/postfix -d
