# alpine-opendkim/rootfs/etc/services.d/opendkim/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/opendkim -x /etc/opendkim/opendkim.conf -u opendkim -P /var/run/opendkim/opendkim.pid -p inet:12345
