# alpine-registry/root/etc/services.d/registry/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid registry /usr/bin/docker-registry /etc/docker-registry/config.yml
